# systemd/consumer-mysql.service
[Unit]
Description=Realtime Ingestion MySQL Consumer
After=network.target rabbitmq-server.service mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
EnvironmentFile=/var/www/app/.env
ExecStartPre=+/bin/mkdir -p /var/www/log/consumer-mysql
ExecStartPre=+/bin/chown -R www-data:www-data /var/www/log/consumer-mysql
ExecStart=/usr/local/bin/consumer-mysql
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=consumer-mysql

[Install]
WantedBy=multi-user.target
