# systemd/consumer-postgres.service
[Unit]
Description=Realtime Ingestion PostgreSQL Consumer
After=network.target rabbitmq-server.service postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
EnvironmentFile=/var/www/app/.env
ExecStartPre=+/bin/mkdir -p /var/www/log/consumer-postgres
ExecStartPre=+/bin/chown -R www-data:www-data /var/www/log/consumer-postgres
ExecStart=/usr/local/bin/consumer-postgres
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=consumer-postgres

[Install]
WantedBy=multi-user.target
